# home/.config/systemd/user/openclaw-gateway.service.d/10-openclaw-secrets.conf
[Service]
EnvironmentFile=/home/goringich/.openclaw/secrets/gateway.env
