# system/etc/systemd/system/bluetooth-self-heal.service
[Unit]
Description=Recover integrated Bluetooth adapter if it becomes unavailable
After=bluetooth.service

[Service]
Type=oneshot
ExecStart=/home/goringich/__home_organized/scripts/bluetooth-self-heal.sh
