# docs/local-ai-stack/snapshot/systemd/codex-agent-orchestrator.service
[Unit]
Description=Run one queued Codex orchestrator task

[Service]
Type=oneshot
ExecStart=%h/.local/bin/codex-agent-run
