# docs/local-ai-stack/snapshot/systemd/ollama-hw-obsidian-log.timer
[Unit]
Description=Run Ollama hardware sampler periodically

[Timer]
OnBootSec=30s
OnUnitActiveSec=60s
AccuracySec=5s
Persistent=true
Unit=ollama-hw-obsidian-log.service

[Install]
WantedBy=timers.target
