# docs/local-ai-stack/snapshot/systemd/openclaw-obsidian-sync.timer
[Unit]
Description=Run OpenClaw to Obsidian sync periodically

[Timer]
OnBootSec=20s
OnUnitActiveSec=2min
Persistent=true
Unit=openclaw-obsidian-sync.service

[Install]
WantedBy=timers.target
