# docs/local-ai-stack/snapshot/systemd/openclaw-obsidian-sync.service
[Unit]
Description=Sync OpenClaw conversations to Obsidian markdown

[Service]
Type=oneshot
ExecStart=/home/goringich/.local/bin/openclaw-obsidian-sync
