# docs/local-ai-stack/snapshot/systemd/openclaw-ollama-sync.service
[Unit]
Description=Sync installed Ollama models into OpenClaw catalogs

[Service]
Type=oneshot
ExecStart=/home/goringich/__home_organized/scripts/openclaw_ollama_sync.py
