# home/.config/systemd/user/swaync.service.d/10-diagnostic-mode.conf
[Unit]
ConditionPathExists=!%h/.local/state/session-startup-guard/disable-swaync
ConditionPathExists=!%h/.local/state/session-startup-guard/disable-swaync-startup
