# home/.config/systemd/user/claude-obsidian-sync.timer
[Unit]
Description=Run Claude to Obsidian sync periodically

[Timer]
OnBootSec=20s
OnUnitActiveSec=2min
Persistent=true
Unit=claude-obsidian-sync.service

[Install]
WantedBy=timers.target
