# home/.config/systemd/user/codex-obsidian-sync.service
[Unit]
Description=Sync Codex conversations to Obsidian markdown
ConditionPathExists=!%h/.local/state/session-startup-guard/disable-nonessential-user-timers

[Service]
Type=oneshot
ExecStart=/home/goringich/.local/bin/codex-obsidian-sync
