# home/.config/systemd/user/codex-obsidian-sync.service
[Unit]
Description=Sync Codex conversations to Obsidian markdown

[Service]
Type=oneshot
ExecStart=/home/goringich/.local/bin/codex-obsidian-sync
