# home/.config/systemd/user/codex-obsidian-sync.timer
[Unit]
Description=Run Codex to Obsidian sync periodically

[Timer]
OnBootSec=20s
OnUnitActiveSec=2min
Persistent=true
Unit=codex-obsidian-sync.service

[Install]
WantedBy=timers.target
