# home/.config/systemd/user/copilot-obsidian-sync.service
[Unit]
Description=Sync Copilot conversations to Obsidian markdown

[Service]
Type=oneshot
ExecStart=/home/goringich/.local/bin/copilot-obsidian-sync
