# home/.config/systemd/user/copilot-obsidian-sync.timer
[Unit]
Description=Run Copilot to Obsidian sync periodically

[Timer]
OnBootSec=20s
OnUnitActiveSec=2min
Persistent=true
Unit=copilot-obsidian-sync.service

[Install]
WantedBy=timers.target
