# home/.config/systemd/user/gpu-watchdog.service
[Unit]
Description=Check GPU health and mirror incidents to Obsidian

[Service]
Type=oneshot
ExecStart=/home/goringich/.config/hypr/scripts/GpuWatchdog.sh
