# home/.config/systemd/user/ollama-hw-obsidian-log.service
[Unit]
Description=Sample Ollama hardware usage and mirror it to Obsidian

[Service]
Type=oneshot
ExecStart=/home/goringich/__home_organized/scripts/ollama_hw_obsidian_logger.py
