# home/.config/systemd/user/openclaw-ollama-sync.timer
[Unit]
Description=Periodically sync installed Ollama models into OpenClaw

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true
Unit=openclaw-ollama-sync.service

[Install]
WantedBy=timers.target
