# home/.config/systemd/user/system-watchdog.service
[Unit]
Description=Check system health and mirror incidents to Obsidian

[Service]
Type=oneshot
ExecStart=/home/goringich/.config/hypr/scripts/SystemWatchdog.sh
